# server/supervisord.conf
[program:vnn-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/vnn/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vnn-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[program:vnn-scheduler]
process_name=%(program_name)s
command=php /var/www/vnn/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vnn-scheduler.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2

[program:vnn-reverb]
process_name=%(program_name)s
command=php /var/www/vnn/artisan reverb:start
autostart=false
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vnn-reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
